# data/server/bundles/teste/teste.conf
[program:teste-george]
command=python -u teste.py
directory=%(ENV_PROJ_DIR)s/src/teste
#stdout_logfile=/home/george/test_supervisor/test_process_output.txt
redirect_stderr=true
autostart=true
autorestart=true
#user=worker
stopsignal=TERM
numprocs=1
process_name=%(program_name)s-%(process_num)s
stdout_logfile=%(ENV_PROJ_DIR)s/data/server/bundles/teste/worker.log
stderr_logfile=%(ENV_PROJ_DIR)s/data/server/bundles/teste/worker-error.log
#environment=PYTHONPATH="$PYTHONPATH:/opt/tornado/"
#environment=PATH="%(ENV_PROJ_DIR)s/lib/pycharm:$PATH"
environment=DEBUG=true


