# data/server/bundles/teste/teste.conf
[program:test_process]
command=python -u teste_supervisor.py
directory=/home/george/test_supervisor
stdout_logfile=/home/george/test_supervisor/test_process_output.txt
redirect_stderr=true
